# hermes-observer/deploy/systemd/linjiang-observer-backend.service.d/10-path.conf
[Service]
Environment=PATH=/home/hermes/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
